# shared/ansible/roles/system/forgejo-runner/templates/forgejo-runner.service.j2
[Unit]
Description=Forgejo Actions Runner
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User={{ superdeploy_user }}
Group={{ superdeploy_group | default(superdeploy_user) }}
WorkingDirectory=/opt/forgejo-runner
ExecStart=/usr/local/bin/forgejo-runner daemon
Restart=always
RestartSec=10
Environment="HOME=/opt/forgejo-runner"

[Install]
WantedBy=multi-user.target
